# context/app/nginx/supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command = /opt/conda/bin/uwsgi --emperor /home/nginx/uwsgi/vassals --master

[program:nginx]
command = /usr/sbin/nginx
